# app/launchd/com.kubota.bias.weekly-recap.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.kubota.bias.weekly-recap</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/kubota/Desktop/bias/app/scripts/weekly_runner.sh</string>
    <string>weekly_recap</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key><integer>1</integer>   <!-- 月曜 -->
    <key>Hour</key><integer>7</integer>
    <key>Minute</key><integer>0</integer>
  </dict>
  <key>WorkingDirectory</key><string>/Users/kubota/Desktop/bias/app</string>
  <key>StandardOutPath</key><string>/Users/kubota/Desktop/bias/app/logs/launchd_weekly-recap.log</string>
  <key>StandardErrorPath</key><string>/Users/kubota/Desktop/bias/app/logs/launchd_weekly-recap.err</string>
</dict>
</plist>
